# rootfs/etc/s6/cron/run
#!/bin/bash
set -eo pipefail

if [ -f ./setup ]
then
  source ./setup
fi

echo "Launching sentry cron..."
su-exec sentry sentry \
  run \
  cron \
  --autoreload
