# rootfs/etc/s6/worker/run
#!/bin/bash
set -eo pipefail

if [ -f ./setup ]
then
  source ./setup
fi

echo "Launching sentry worker..."
su-exec sentry sentry \
  run \
  worker \
  --autoreload
